# data/indicator-session.service.in
[Unit]
Description=Indicator Session Service
PartOf=graphical-session.target
PartOf=indicators.target

[Service]
ExecStart=@pkglibexecdir@/indicator-session-service
Restart=on-failure

[Install]
WantedBy=indicators.target
